# deploy/skywatch.service
[Unit]
Description=Skywatch ADS-B Tracker
After=network.target dump1090.service
Wants=dump1090.service

[Service]
Type=simple
User=skywatch
Group=skywatch
ExecStart=/opt/skywatch/adsb-tracker -config /etc/skywatch/config.json
WorkingDirectory=/opt/skywatch
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

